# crontab
* * * * *   updatedb -l 0 -o ~/projects/ldb -U ~/projects --prunenames=".git .project .projects target .class WEB-INF/classes .sass-cache .pyc"
